# other/crontab
# m h  dom mon dow   command
00,15,30,45 * * * * /srv/check_wakebird